# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:flask]
command=python -c "from backend.app import create_app; from backend.app.extensions import socketio; app = create_app(); socketio.run(app, host='127.0.0.1', port=5000, allow_unsafe_werkzeug=True)"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/flask.log
stderr_logfile=/var/log/supervisor/flask_err.log
environment=PYTHONPATH="/app",FLASK_DEBUG="0"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log

; 监听重启标记文件，支持一键升级后自动重启 Flask
[program:restart-watcher]
command=bash -c "while true; do if [ -f /app/.restart-requested ]; then rm -f /app/.restart-requested; supervisorctl restart flask; fi; sleep 5; done"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/restart-watcher.log
stderr_logfile=/dev/null
